# tests/docker/supervisor.conf
# Copyright 2016 Workiva Inc.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

[supervisord]
nodaemon=true

[program:memcached1]
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
command=memcached -u root
priority: 1

[program:memcached2]
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
command=memcached -u root -p 22122
priority: 1

[program:dev_lambda.py1]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "sleep 5; python tools/dev_lambda.py --kinesis_stream_arn=PRIMARY_STREAM_SOURCE --run_kinesis_lambda=1 --sleep_time=0.01"
priority: 9

[program:dev_lambda.py2]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "sleep 5; python tools/dev_lambda.py --kinesis_stream_arn=SECONDARY_STREAM_SOURCE --run_kinesis_lambda=1 --sleep_time=0.01"
priority: 9

[program:dev_lambda.py3]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "sleep 5; python tools/dev_lambda.py --kinesis_stream_arn=PRIMARY_RETRY_SOURCE --run_kinesis_lambda=1 --sleep_time=0.01"
priority: 9

[program:dev_lambda.py4]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "sleep 5; python tools/dev_lambda.py --kinesis_stream_arn=SECONDARY_RETRY_SOURCE --run_kinesis_lambda=1 --sleep_time=0.01"
priority: 9

[program:kinesalite1]
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
command=kinesalite --port 1111
priority: 5

[program:kinesalite2]
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
command=kinesalite --port 2222
priority: 5

[program:kinesalite3]
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
command=kinesalite --port 3333
priority: 5

[program:kinesalite4]
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
command=kinesalite --port 4444
priority: 5

[program:setup]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "python tools/create_kinesis_stream.py --kinesis_stream_arn=PRIMARY_STREAM_SOURCE;
  python tools/create_kinesis_stream.py --kinesis_stream_arn=SECONDARY_STREAM_SOURCE;
  python tools/create_kinesis_stream.py --kinesis_stream_arn=PRIMARY_RETRY_SOURCE;
  python tools/create_kinesis_stream.py --kinesis_stream_arn=SECONDARY_RETRY_SOURCE;
  sleep 3600"
priority: 8

[program:start]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "sleep 10; python tools/start_state_machine.py --machine_name=tracer; sleep 3600"
priority: 100