# tests/functional/supervisor.conf
# Copyright 2016-2020 Workiva Inc.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

[supervisord]
nodaemon=true

[program:dev_lambda_stream_primary]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "while [ ! -f /tmp/setup ]; do sleep 1; done; %(ENV_PYTHON_BIN)s tools/dev_lambda.py --sqs_queue_arn=PRIMARY_STREAM_SOURCE --run_sqs_lambda=1 --sleep_time=0.05"
priority: 9

[program:dev_lambda_stream_secondary]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "while [ ! -f /tmp/setup ]; do sleep 1; done; %(ENV_PYTHON_BIN)s tools/dev_lambda.py --kinesis_stream_arn=SECONDARY_STREAM_SOURCE --run_kinesis_lambda=1 --sleep_time=0.05"
priority: 9

[program:dev_lambda_retry_primary]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "while [ ! -f /tmp/setup ]; do sleep 1; done; %(ENV_PYTHON_BIN)s tools/dev_lambda.py --sqs_queue_arn=PRIMARY_RETRY_SOURCE --run_sqs_lambda=1 --sleep_time=0.05"
priority: 9

[program:dev_lambda_retry_secondary]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "while [ ! -f /tmp/setup ]; do sleep 1; done; %(ENV_PYTHON_BIN)s tools/dev_lambda.py --kinesis_stream_arn=SECONDARY_RETRY_SOURCE --run_kinesis_lambda=1 --sleep_time=0.05"
priority: 9

[program:setup]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=.
command=/bin/sh -c "
  sleep 30;
  %(ENV_PYTHON_BIN)s tools/create_sqs_queue.py --sqs_queue_arn=PRIMARY_STREAM_SOURCE || true;
  %(ENV_PYTHON_BIN)s tools/create_kinesis_stream.py --kinesis_stream_arn=SECONDARY_STREAM_SOURCE || true;
  %(ENV_PYTHON_BIN)s tools/create_sqs_queue.py --sqs_queue_arn=PRIMARY_RETRY_SOURCE || true;
  %(ENV_PYTHON_BIN)s tools/create_kinesis_stream.py --kinesis_stream_arn=SECONDARY_RETRY_SOURCE || true;
  %(ENV_PYTHON_BIN)s tools/create_dynamodb_table.py --dynamodb_table_arn=PRIMARY_CHECKPOINT_SOURCE || true;
  %(ENV_PYTHON_BIN)s tools/create_dynamodb_table.py --dynamodb_table_arn=PRIMARY_ENVIRONMENT_SOURCE || true;
  touch /tmp/setup;
  sleep 3600"
priority: 8
